# .github/workflows/packaging/deb/postinst
#!/bin/bash

SERVICE="/opt/linupdate/service.py"

# Restore configuration files if exists
if [ -f "/tmp/linupdate.yml.debsave" ];then
    rm -f /etc/linupdate/linupdate.yml
    mv /tmp/linupdate.yml.debsave /etc/linupdate/linupdate.yml
fi
if [ -f "/tmp/update.yml.debsave" ];then
    rm -f /etc/linupdate/update.yml
    mv /tmp/update.yml.debsave /etc/linupdate/update.yml
fi

# If no configuration files exists, copy default
if [ ! -f "/etc/linupdate/linupdate.yml" ];then
    cp /opt/linupdate/templates/linupdate.template.yml /etc/linupdate/linupdate.yml
fi
if [ ! -f "/etc/linupdate/update.yml" ];then
    cp /opt/linupdate/templates/update.template.yml /etc/linupdate/update.yml
fi

# Create a symlink to main script
ln -sf /opt/linupdate/linupdate.py /usr/bin/linupdate

# Install en_US.UTF-8 locale if not present
if ! locale -a | grep -q "en_US.UTF-8";then
    apt-get install locales-all -y > /dev/null
fi

# Set permissions
chmod 750 /etc/linupdate
chmod 750 /opt/linupdate

# Only if systemd is installed (not the case on github runners)
if [ -f "/usr/bin/systemctl" ];then
    # Copy systemd unit file if not exists
    if [ ! -f "/lib/systemd/system/linupdate.service" ];then
        cp /opt/linupdate/service/linupdate.systemd.template /etc/systemd/system/linupdate.service
    fi

    # Enable service script by creating a symlink if not exists
    if [ ! -L "/etc/systemd/system/linupdate.service" ];then
        # Delete file in case it is a file and not a symlink
        rm -f /etc/systemd/system/linupdate.service
        ln -sf /lib/systemd/system/linupdate.service /etc/systemd/system/linupdate.service
    fi

    chmod 550 "$SERVICE"
    chown root:root "$SERVICE"

    # Clean old directories
    rm /opt/linupdate/agents-enabled -rf
    rm /opt/linupdate/mods-enabled -rf
    rm /opt/linupdate/mods-available -rf

    /usr/bin/systemctl enable linupdate
    /usr/bin/systemctl --quiet daemon-reload

    # Start service
    if /usr/bin/systemctl is-active --quiet linupdate;then
        /usr/bin/systemctl restart --quiet linupdate
    else
        /usr/bin/systemctl start --quiet linupdate
    fi
fi